# initConf/celerybeat.ini
[program:celery_beat]
command=celery -A NetOpsBase beat --loglevel=info --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/Users/congxingwang/pythoncode/NetOpsBase
stdout_logfile=/Users/congxingwang/pythoncode/NetOpsBase/logs/celerybeat_out.log
stderr_logfile=/Users/congxingwang/pythoncode/NetOpsBase/logs/celerybeat_err.log
autorestart=true
redirect_stderr=true